# upstart/switchboard.conf
description "Switchboard redirect server"
author      "Jamie Liu <jliu@nixprime.com>"

start on net-device-up IFACE!=lo
stop on runlevel [!2345]

setuid nobody

exec /usr/local/bin/switchboard
